# installers/legacy/network-optimize.service
[Unit]
Description=Network Optimization and Load Balancing Service
After=network-online.target NetworkManager.service systemd-networkd.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/network-optimize.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

# Restart on failure
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
