# installers/legacy/network-optimize-periodic.service
[Unit]
Description=Network Re-optimization (Periodic)
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/network-optimize.sh
StandardOutput=journal
StandardError=journal
